# infra/systemd/eyye-ingest.timer
# file: infra/systemd/eyye-ingest.timer
[Unit]
Description=Run EYYE unified ingest (Telegram + Wikipedia) hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
